# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# start DLNA services
if [ -e /etc/dlna.sh ]; then
/etc/dlna.sh &
fi
exit 0
